# scripts/sync-packages.service
[Unit]
Description=Sync Arch packages from Git repository

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sync-packages.sh
